# Miapp.service
[Unit]
Description=gunicorn daemon
After=networking.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/BDFeslic
Environment="PATH=cd"
ExecStart=sudo sh /home/ubuntu/BDFeslic/script.sh
[Install]
WantedBy=multi-user.target
